# server/systemd/user/server-auto-update.service
[Unit]
Description=Unattended server maintenance (npm updates)

[Service]
Type=oneshot
ExecStart=%h/.local/bin/server-auto-update
